# ansible/roles/deploy_prometheus/templates/node-exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
User=root

[Service]
ExecStart=/opt/prometheus/node_exporter-1.8.1.linux-amd64/node_exporter
